# porteus/porteus/rootcopy/etc/rc.d/rc.local
#!/bin/sh

sudo "echo performance | tee /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor"
sudo "echo 1 | tee /sys/devices/system/cpu/cpufreq/boost"

chmod -R +x /root/tools
